# yo_trigger/pkgroot/Library/LaunchAgents/com.github.groob.yo-trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Applications/Server.app/Contents/ServerRoot/usr/bin:/Applications/Server.app/Contents/ServerRoot/usr/sbin:/usr/local/go/bin:/usr/local/munki:/usr/local/sbin</string>
	</dict>
	<key>Label</key>
	<string>com.github.groob.yo-trigger</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/yo-trigger</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/Library/Management/Triggers/yo</string>
	</array>
</dict>
</plist>
