# logstash_forwarder/pkgroot/Library/LaunchDaemons/com.elastic.logstash-forwarder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.elastic.logstash-forwarder</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/logstash-forwarder</string>
		<string>-config</string>
		<string>/etc/logstash-forwarder/config</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/log/logstash-forwarder.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/logstash-forwarder.log</string>
	<key>WorkingDirectory</key>
	<string>/var/run</string>
</dict>
</plist>
